# DEPRECATED/laravel-boost-mcp/laravel-boost-mcp.conf
# /etc/supervisor/conf.d/mcp.conf

[program:mcp]
command=laravel-boost-mcp
autostart=true
autorestart=true
